# ansible-runner/rootfs/etc/services.d/example/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the example service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

# Add your code here

# Declare variables
declare message
declare url
declare repo
declare branch

## Get the 'message' key from the user config options.
message=$(bashio::config 'message')
url=$(bashio::config 'repo')
repo=$(echo $URL | sed -E 's@https:\/\/github\.com\/[^\/]+\/@@' | sed -E 's@.git$@@')
branch=$(bashio::config 'branch')

## Run your program
git clone https://username:password@host

exec rm -rf $repo && git clone $url && cd $repo && git checkout -b ${branch:='main'} && ansible-playbook site.yml

COPY data/run.sh /

